# hosting-services/gcp/sepal/profile.d/multitail.sh
if [ -x /bin/multitail ]; then
    alias mt='multitail -CT ANSI'
fi
