# sepal-sandbox/docker/sandbox/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:install_oft]
command=/install_oft.sh