# modules/rsyslog/docker/config/supervisord.conf
[supervisord]
nodaemon=true
loglevel=info

[program:rsyslog]
command=/usr/local/bin/start.sh -n
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
