# haproxy/docker/config/supervisord.conf
[supervisord]
nodaemon=true

[program:haproxy]
command=haproxy -f /usr/local/etc/haproxy/haproxy.cfg
redirect_stderr=true
redirect_stdin=true
startretries=0


[program:rsyslog]
command=/etc/init.d/rsyslog start
redirect_stderr=true
redirect_stdin=true
startretries=0