# ssh-gateway/docker/config/supervisord.conf
[supervisord]
nodaemon=true

[program:incron]
command=/etc/init.d/incron start
redirect_stderr=true
redirect_stdin=true
autorestart=false
startretries=0

[program:sshd]
command=/usr/sbin/sshd -D


