# hosting-services/gcp/sepal/gpu/init-gpu-drivers.service
[Unit]
Description=Init GPU drivers

[Service]
Type=oneshot
ExecStart=/usr/local/bin/init-gpu-drivers.sh

[Install]
WantedBy=multi-user.target
