# sepal-sandbox/docker/mysql/config/supervisord.conf
[supervisord]
nodaemon=true

[program:init]
command=/init_container.sh

[program:mysql]
command=/entrypoint.sh mysqld
redirect_stderr=true
redirect_stdin=true
priority=1






